# sites/sanexio-portal/install/de.sanexio.archiv-sync-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.sanexio.archiv-sync-api</string>

    <key>WorkingDirectory</key>
    <string>${HOME}/Cortex/projects/Cortex-Web/sites/sanexio-portal</string>

    <key>ProgramArguments</key>
    <array>
        <string>__NODE_BIN__</string>
        <string>server/archiv-sync-api.mjs</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__NODE_DIR__:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>ARCHIV_SYNC_HOST</key>
        <string>127.0.0.1</string>
        <key>ARCHIV_SYNC_PORT</key>
        <string>5177</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>${HOME}/Library/Logs/cortex-services/archiv-sync-api.log</string>

    <key>StandardErrorPath</key>
    <string>${HOME}/Library/Logs/cortex-services/archiv-sync-api.err.log</string>

    <key>ThrottleInterval</key>
    <integer>15</integer>
</dict>
</plist>
