# sites/sanexio-portal/install/de.sanexio.workforce-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.sanexio.workforce-api</string>

    <key>WorkingDirectory</key>
    <string>${HOME}/Cortex/projects/Cortex-Web/sites/workforce-time</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/node@22/bin/node</string>
        <string>server/api.js</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/opt/node@22/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>NODE_ENV</key>
        <string>development</string>
        <key>CORTEX_TENANT_DIR</key>
        <string>${HOME}/Cortex/projects/Sanexio-Tenant</string>
        <key>ARBEITSZEITEN_API_HOST</key>
        <string>127.0.0.1</string>
        <key>ARBEITSZEITEN_API_PORT</key>
        <string>5175</string>
        <key>ARBEITSZEITEN_EXTRA_ORIGINS</key>
        <string>https://${LOCAL_WP_HOST}</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>${HOME}/Library/Logs/cortex-services/workforce-api.log</string>

    <key>StandardErrorPath</key>
    <string>${HOME}/Library/Logs/cortex-services/workforce-api.err.log</string>

    <key>ThrottleInterval</key>
    <integer>15</integer>
</dict>
</plist>
